# ollama.service
[Unit]
Description=Ollama LLM Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
